# bin/invm/charmander-scheduler.conf
description "charmander scheduler"

start on stopped rc RUNLEVEL=[2345]
respawn

exec /usr/local/bin/charmander-scheduler -local-ip="172.31.1.11" -master="172.31.1.11:5050"