# vector/pcp/supervisord.conf
[supervisord]
nodaemon=true

[program:pmcd]
command=pmcd start

[program:pmwebd]
command=pmwebd start
